# deploy/oracle/dealbot-post.service
[Unit]
Description=Deal bot - post run (drain 1 queued deal to the channel)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/.venv/bin/python -m app.main --post-queue 1
TimeoutStartSec=300
